# systemd/user/hermesops-orchestrator.service
[Unit]
Description=HermesOps persistent multi-task DAG orchestrator
Documentation=file:///opt/docker/hermesops/repo/docs/ORCHESTRATION.md
After=hermesops-supervisor.service
Wants=hermesops-supervisor.service
StartLimitIntervalSec=60
StartLimitBurst=6

[Service]
Type=simple
WorkingDirectory=/opt/docker/hermesops/repo
Environment=HERMESOPS_ROOT=/opt/docker/hermesops
ExecStart=/usr/bin/python3 /opt/docker/hermesops/repo/scripts/hermesops-orchestrator.py daemon
Restart=always
RestartSec=5s
TimeoutStopSec=120s
UMask=0027

NoNewPrivileges=true
PrivateTmp=true
RestrictSUIDSGID=true
LockPersonality=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
