# systemd/user/hermesops-controller-api.service
[Unit]
Description=HermesOps local read-only Controller API
Documentation=file:///opt/docker/hermesops/repo/docs/milestones/2C_CONTROLLER_API_SERVICE.md
StartLimitIntervalSec=60
StartLimitBurst=6

[Service]
Type=simple
WorkingDirectory=/opt/docker/hermesops/repo
Environment=HERMESOPS_ROOT=/opt/docker/hermesops
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONPATH=/opt/docker/hermesops/repo
ExecStartPre=/usr/bin/python3 /opt/docker/hermesops/repo/scripts/hermesops-controller-session.py check
ExecStartPre=/usr/bin/python3 /opt/docker/hermesops/repo/scripts/hermesops-controller-api.py check --host 127.0.0.1 --port 8765
ExecStart=/usr/bin/python3 /opt/docker/hermesops/repo/scripts/hermesops-controller-api.py serve --host 127.0.0.1 --port 8765 --log-level INFO
ExecStartPost=/usr/bin/python3 /opt/docker/hermesops/repo/scripts/hermesops-controller-probe.py --base-url http://127.0.0.1:8765 --wait-seconds 20
Restart=on-failure
RestartSec=5s
TimeoutStartSec=40s
TimeoutStopSec=20s
KillSignal=SIGTERM
KillMode=mixed
UMask=0077

# Portable hardening for a non-root systemd --user manager.
# Capability-bounding and privileged mount-namespace directives are excluded
# because they fail on supported hosts without private user namespaces.
NoNewPrivileges=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
