# systemd/user/hermesops-console.service
[Unit]
Description=HermesOps local Console foundation
Documentation=file:///opt/docker/hermesops/repo/docs/milestones/2P_CONSOLE_WEB_FOUNDATION.md
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=6

[Service]
Type=simple
WorkingDirectory=/opt/docker/hermesops/repo
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStartPre=/usr/bin/python3 /opt/docker/hermesops/repo/scripts/hermesops-console.py check --root /opt/docker/hermesops/repo/console/dist --host 127.0.0.1 --port 8788
ExecStart=/usr/bin/python3 /opt/docker/hermesops/repo/scripts/hermesops-console.py serve --root /opt/docker/hermesops/repo/console/dist --host 127.0.0.1 --port 8788 --max-connections 16
ExecStartPost=/usr/bin/python3 /opt/docker/hermesops/repo/scripts/hermesops-console-probe.py --base-url http://127.0.0.1:8788 --wait-seconds 20
Restart=on-failure
RestartSec=5s
TimeoutStartSec=35s
TimeoutStopSec=20s
KillSignal=SIGTERM
KillMode=mixed
UMask=0077

NoNewPrivileges=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
