# systemd/user/hermesops-supervisor.service
[Unit]
Description=HermesOps automatic Supervisor and Watchdog
Documentation=file:///opt/docker/hermesops/repo/docs/SUPERVISOR.md
After=default.target

[Service]
Type=simple
WorkingDirectory=/opt/docker/hermesops/repo
Environment=HERMESOPS_ROOT=/opt/docker/hermesops
ExecStart=/usr/bin/python3 /opt/docker/hermesops/repo/scripts/hermesops-supervisor.py run
Restart=always
RestartSec=5s
TimeoutStopSec=30s
UMask=0027

NoNewPrivileges=true
PrivateTmp=true
RestrictSUIDSGID=true
LockPersonality=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
