# systemd/user/hermesops-notifier.service
[Unit]
Description=HermesOps durable notification outbox
StartLimitIntervalSec=60
StartLimitBurst=6

[Service]
Type=simple
WorkingDirectory=/opt/docker/hermesops/repo
ExecStart=/usr/bin/python3 /opt/docker/hermesops/repo/scripts/hermesops-notifier.py daemon
Restart=always
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=default.target
